# deb/supervisor-config/supervisor-hhvm.conf
[program:hhvm]
command=/scripts/supervisor-hhvm.sh
numprocs=1
autostart=true
autorestart=true