# deb/supervisord-hhvm.conf
[program:hhvm]
command=/supervisord-hhvm.sh
numprocs=1
autostart=true
autorestart=true